# supervisord.conf
[supervisord]
nodaemon=false

[program:xvfb]
command=/bin/sh -c "Xvfb :100 -screen 0 1024x768x24"

[program:sshd]
command=/usr/sbin/sshd -D

